# fs/alpine/etc/profile.d/05-color.sh
COLBLACK=$(tput setaf 0)
COLMAROON=$(tput setaf 1)
COLGREEN=$(tput setaf 2)
COLOLIVE=$(tput setaf 3)
COLNAVYBLUE=$(tput setaf 4)
COLPURPLE=$(tput setaf 5)
COLTEAL=$(tput setaf 6)
COLSILVER=$(tput setaf 7)
COLGREY=$(tput setaf 8)
COLRED=$(tput setaf 9)
COLLIME=$(tput setaf 10)
COLYELLOW=$(tput setaf 11)
COLBLUE=$(tput setaf 12)
COLFUCHSIA=$(tput setaf 13)
COLAQUA=$(tput setaf 14)
COLWHITE=$(tput setaf 15)
COLNAVYBLUE=$(tput setaf 17)
COLDARKBLUE=$(tput setaf 18)
COLBLUE3=$(tput setaf 19)
COLBLUE3=$(tput setaf 20)
COLBLUE1=$(tput setaf 21)
COLDARKGREEN=$(tput setaf 22)
COLDEEPSKYBLUE4=$(tput setaf 23)
COLDEEPSKYBLUE4=$(tput setaf 24)
COLDEEPSKYBLUE4=$(tput setaf 25)
COLDODGERBLUE3=$(tput setaf 26)
COLDODGERBLUE2=$(tput setaf 27)
COLGREEN4=$(tput setaf 28)
COLSPRINGGREEN4=$(tput setaf 29)
COLTURQUOISE4=$(tput setaf 30)
COLDEEPSKYBLUE3=$(tput setaf 31)
COLDEEPSKYBLUE3=$(tput setaf 32)
COLDODGERBLUE1=$(tput setaf 33)
COLGREEN3=$(tput setaf 34)
COLSPRINGGREEN3=$(tput setaf 35)
COLDARKCYAN=$(tput setaf 36)
COLLIGHTSEAGREEN=$(tput setaf 37)
COLDEEPSKYBLUE2=$(tput setaf 38)
COLDEEPSKYBLUE1=$(tput setaf 39)
COLGREEN3=$(tput setaf 40)
COLSPRINGGREEN3=$(tput setaf 41)
COLSPRINGGREEN2=$(tput setaf 42)
COLCYAN3=$(tput setaf 43)
COLDARKTURQUOISE=$(tput setaf 44)
COLTURQUOISE2=$(tput setaf 45)
COLGREEN1=$(tput setaf 46)
COLSPRINGGREEN2=$(tput setaf 47)
COLSPRINGGREEN1=$(tput setaf 48)
COLMIDSPRINGGREEN=$(tput setaf 49)
COLCYAN2=$(tput setaf 50)
COLCYAN1=$(tput setaf 51)
COLDARKRED=$(tput setaf 52)
COLDEEPPINK4=$(tput setaf 53)
COLPURPLE4=$(tput setaf 54)
COLPURPLE4=$(tput setaf 55)
COLPURPLE3=$(tput setaf 56)
COLBLUEVIOLET=$(tput setaf 57)
COLORANGE4=$(tput setaf 58)
COLGREY37=$(tput setaf 59)
COLMIDPURPLE4=$(tput setaf 60)
COLSLATEBLUE3=$(tput setaf 61)
COLSLATEBLUE3=$(tput setaf 62)
COLROYALBLUE1=$(tput setaf 63)
COLCHARTREUSE4=$(tput setaf 64)
COLDARKSEAGREEN4=$(tput setaf 65)
COLPALETURQUOISE4=$(tput setaf 66)
COLSTEELBLUE=$(tput setaf 67)
COLSTEELBLUE3=$(tput setaf 68)
COLCORNFLOWERBLUE=$(tput setaf 69)
COLCHARTREUSE3=$(tput setaf 70)
COLDARKSEAGREEN4=$(tput setaf 71)
COLCADETBLUE=$(tput setaf 72)
COLCADETBLUE=$(tput setaf 73)
COLSKYBLUE3=$(tput setaf 74)
COLSTEELBLUE1=$(tput setaf 75)
COLCHARTREUSE3=$(tput setaf 76)
COLPALEGREEN3=$(tput setaf 77)
COLSEAGREEN3=$(tput setaf 78)
COLAQUAMARINE3=$(tput setaf 79)
COLMIDTURQUOISE=$(tput setaf 80)
COLSTEELBLUE1=$(tput setaf 81)
COLCHARTREUSE2=$(tput setaf 82)
COLSEAGREEN2=$(tput setaf 83)
COLSEAGREEN1=$(tput setaf 84)
COLSEAGREEN1=$(tput setaf 85)
COLAQUAMARINE1=$(tput setaf 86)
COLDARKSLATEGRAY2=$(tput setaf 87)
COLDARKRED=$(tput setaf 88)
COLDEEPPINK4=$(tput setaf 89)
COLDARKMAGENTA=$(tput setaf 90)
COLDARKMAGENTA=$(tput setaf 91)
COLDARKVIOLET=$(tput setaf 92)
COLPURPLE=$(tput setaf 93)
COLORANGE4=$(tput setaf 94)
COLLIGHTPINK4=$(tput setaf 95)
COLPLUM4=$(tput setaf 96)
COLMIDPURPLE3=$(tput setaf 97)
COLMIDPURPLE3=$(tput setaf 98)
COLSLATEBLUE1=$(tput setaf 99)
COLYELLOW4=$(tput setaf 100)
COLWHEAT4=$(tput setaf 101)
COLGREY53=$(tput setaf 102)
COLLIGHTSLATEGREY=$(tput setaf 103)
COLMIDPURPLE=$(tput setaf 104)
COLLIGHTSLATEBLUE=$(tput setaf 105)
COLYELLOW4=$(tput setaf 106)
COLDARKOLIVEGREEN3=$(tput setaf 107)
COLDARKSEAGREEN=$(tput setaf 108)
COLLIGHTSKYBLUE3=$(tput setaf 109)
COLLIGHTSKYBLUE3=$(tput setaf 110)
COLSKYBLUE2=$(tput setaf 111)
COLCHARTREUSE2=$(tput setaf 112)
COLDARKOLIVEGREEN3=$(tput setaf 113)
COLPALEGREEN3=$(tput setaf 114)
COLDARKSEAGREEN3=$(tput setaf 115)
COLDARKSLATEGRAY3=$(tput setaf 116)
COLSKYBLUE1=$(tput setaf 117)
COLCHARTREUSE1=$(tput setaf 118)
COLLIGHTGREEN=$(tput setaf 119)
COLLIGHTGREEN=$(tput setaf 120)
COLPALEGREEN1=$(tput setaf 121)
COLAQUAMARINE1=$(tput setaf 122)
COLDARKSLATEGRAY1=$(tput setaf 123)
COLRED3=$(tput setaf 124)
COLDEEPPINK4=$(tput setaf 125)
COLMIDVIOLETRED=$(tput setaf 126)
COLMAGENTA3=$(tput setaf 127)
COLDARKVIOLET=$(tput setaf 128)
COLPURPLE=$(tput setaf 129)
COLDARKORANGE3=$(tput setaf 130)
COLINDIANRED=$(tput setaf 131)
COLHOTPINK3=$(tput setaf 132)
COLMIDORCHID3=$(tput setaf 133)
COLMIDORCHID=$(tput setaf 134)
COLMIDPURPLE2=$(tput setaf 135)
COLDARKGOLDENROD=$(tput setaf 136)
COLLIGHTSALMON3=$(tput setaf 137)
COLROSYBROWN=$(tput setaf 138)
COLGREY63=$(tput setaf 139)
COLMIDPURPLE2=$(tput setaf 140)
COLMIDPURPLE1=$(tput setaf 141)
COLGOLD3=$(tput setaf 142)
COLDARKKHAKI=$(tput setaf 143)
COLNAVAJOWHITE3=$(tput setaf 144)
COLGREY69=$(tput setaf 145)
COLLIGHTSTEELBLUE3=$(tput setaf 146)
COLLIGHTSTEELBLUE=$(tput setaf 147)
COLYELLOW3=$(tput setaf 148)
COLDARKOLIVEGREEN3=$(tput setaf 149)
COLDARKSEAGREEN3=$(tput setaf 150)
COLDARKSEAGREEN2=$(tput setaf 151)
COLLIGHTCYAN3=$(tput setaf 152)
COLLIGHTSKYBLUE1=$(tput setaf 153)
COLGREENYELLOW=$(tput setaf 154)
COLDARKOLIVEGREEN2=$(tput setaf 155)
COLPALEGREEN1=$(tput setaf 156)
COLDARKSEAGREEN2=$(tput setaf 157)
COLDARKSEAGREEN1=$(tput setaf 158)
COLPALETURQUOISE1=$(tput setaf 159)
COLRED3=$(tput setaf 160)
COLDEEPPINK3=$(tput setaf 161)
COLDEEPPINK3=$(tput setaf 162)
COLMAGENTA3=$(tput setaf 163)
COLMAGENTA3=$(tput setaf 164)
COLMAGENTA2=$(tput setaf 165)
COLDARKORANGE3=$(tput setaf 166)
COLINDIANRED=$(tput setaf 167)
COLHOTPINK3=$(tput setaf 168)
COLHOTPINK2=$(tput setaf 169)
COLORCHID=$(tput setaf 170)
COLMIDORCHID1=$(tput setaf 171)
COLORANGE3=$(tput setaf 172)
COLLIGHTSALMON3=$(tput setaf 173)
COLLIGHTPINK3=$(tput setaf 174)
COLPINK3=$(tput setaf 175)
COLPLUM3=$(tput setaf 176)
COLVIOLET=$(tput setaf 177)
COLGOLD3=$(tput setaf 178)
COLLIGHTGOLDENROD3=$(tput setaf 179)
COLTAN=$(tput setaf 180)
COLMISTYROSE3=$(tput setaf 181)
COLTHISTLE3=$(tput setaf 182)
COLPLUM2=$(tput setaf 183)
COLYELLOW3=$(tput setaf 184)
COLKHAKI3=$(tput setaf 185)
COLLIGHTGOLDENROD2=$(tput setaf 186)
COLLIGHTYELLOW3=$(tput setaf 187)
COLGREY84=$(tput setaf 188)
COLLIGHTSTEELBLUE1=$(tput setaf 189)
COLYELLOW2=$(tput setaf 190)
COLDARKOLIVEGREEN1=$(tput setaf 191)
COLDARKOLIVEGREEN1=$(tput setaf 192)
COLDARKSEAGREEN1=$(tput setaf 193)
COLHONEYDEW2=$(tput setaf 194)
COLLIGHTCYAN1=$(tput setaf 195)
COLRED1=$(tput setaf 196)
COLDEEPPINK2=$(tput setaf 197)
COLDEEPPINK1=$(tput setaf 198)
COLDEEPPINK1=$(tput setaf 199)
COLMAGENTA2=$(tput setaf 200)
COLMAGENTA1=$(tput setaf 201)
COLORANGERED1=$(tput setaf 202)
COLINDIANRED1=$(tput setaf 203)
COLINDIANRED1=$(tput setaf 204)
COLHOTPINK=$(tput setaf 205)
COLHOTPINK=$(tput setaf 206)
COLMIDORCHID1=$(tput setaf 207)
COLDARKORANGE=$(tput setaf 208)
COLSALMON1=$(tput setaf 209)
COLLIGHTCORAL=$(tput setaf 210)
COLPALEVIOLETRED1=$(tput setaf 211)
COLORCHID2=$(tput setaf 212)
COLORCHID1=$(tput setaf 213)
COLORANGE1=$(tput setaf 214)
COLSANDYBROWN=$(tput setaf 215)
COLLIGHTSALMON1=$(tput setaf 216)
COLLIGHTPINK1=$(tput setaf 217)
COLPINK1=$(tput setaf 218)
COLPLUM1=$(tput setaf 219)
COLGOLD1=$(tput setaf 220)
COLLIGHTGOLDENROD2=$(tput setaf 221)
COLLIGHTGOLDENROD2=$(tput setaf 222)
COLNAVAJOWHITE1=$(tput setaf 223)
COLMISTYROSE1=$(tput setaf 224)
COLTHISTLE1=$(tput setaf 225)
COLYELLOW1=$(tput setaf 226)
COLLIGHTGOLDENROD1=$(tput setaf 227)
COLKHAKI1=$(tput setaf 228)
COLWHEAT1=$(tput setaf 229)
COLCORNSILK1=$(tput setaf 230)
COLGREY100=$(tput setaf 231)
COLGREY3=$(tput setaf 232)
COLGREY7=$(tput setaf 233)
COLGREY11=$(tput setaf 234)
COLGREY15=$(tput setaf 235)
COLGREY19=$(tput setaf 236)
COLGREY23=$(tput setaf 237)
COLGREY27=$(tput setaf 238)
COLGREY30=$(tput setaf 239)
COLGREY35=$(tput setaf 240)
COLGREY39=$(tput setaf 241)
COLGREY42=$(tput setaf 242)
COLGREY46=$(tput setaf 243)
COLGREY50=$(tput setaf 244)
COLGREY54=$(tput setaf 245)
COLGREY58=$(tput setaf 246)
COLGREY62=$(tput setaf 247)
COLGREY66=$(tput setaf 248)
COLGREY70=$(tput setaf 249)
COLGREY74=$(tput setaf 250)
COLGREY78=$(tput setaf 251)
COLGREY82=$(tput setaf 252)
COLGREY85=$(tput setaf 253)
COLGREY89=$(tput setaf 254)
COLGREY93=$(tput setaf 255)
